# vendor_docker/server/cron/crontab
* * * * * cd /var/www/html/public && php artisan schedule:run >> /dev/stdout 2>&1
